# systemd/pi-provision.service
[Unit]
Description=Pi Provisioning (one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=jr
WorkingDirectory=/home/jr/pi-provision
ExecStart=/bin/bash /home/jr/pi-provision/provision.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
